# tw5.service
[Unit]
Description=TiddlyWiki5 Service

[Service]
WorkingDirectory=/home/oeyoews/REPOS/TiddlyWiki5
ExecStart=/usr/bin/tiddlywiki --listen port=8090 anon-username=oeyoews use-browser-cache=yes gzip=yes
